# platform/filestore/config/blinker-filestore.conf
[program:blinker-filestore]
command = blinker-filestore
user = blinker
environment = BLINKER_FILESTORE_CONFIGFILE="/etc/blinker/filestore.yml"
startsecs = 10