# platform/web/config/blinker-web.conf
[program:blinker-web]
command = blinker-web
user = blinker
environment = BLINKER_WEB_CONFIGFILE="/etc/blinker/web.yml"
startsecs = 10